# s3_resources/dns_updater/crontab
@reboot valheim python3 /opt/tools/dns_updater/dns_updater.py