# s3_resources/metrics/crontab
*/5 * * * * valheim /bin/sh /opt/tools/metrics/collect.sh