# s3_resources/dns_update/crontab
@reboot minecraft python3 /opt/minecraft_aws_tools/dns_updater/dns_updater.py